# supervisord.conf
[unix_http_server]
file=/tmp/senpy.sock                       ; path to your socket file

[supervisord]
logfile = %(here)s/logs/supervisor.log
childlogdir = %(here)s/logs/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
logfile = %(here)s/logs/supervisorctl.log
serverurl=unix:///tmp/senpy.sock         ; use a unix:// URL  for a unix socket


[program:senpy]
command = venv/bin/gunicorn -w 1 -b 0.0.0.0:6666 --log-file %(here)s/logs/gunicorn.log --access-logfile - app:app
directory = %(here)s
environment = PATH=%(here)s/venv/bin/
logfile = %(here)s/logs/senpy.log
